# scripts/podman/openclaw.container.in
# elysiaclaw gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into elysiaclaw's ~/.config/containers/systemd/
# {{elysiaclaw_HOME}} is replaced at install time.

[Unit]
Description=elysiaclaw gateway (rootless Podman)

[Container]
Image=elysiaclaw:local
ContainerName=elysiaclaw
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{elysiaclaw_HOME}}/.elysiaclaw:/home/node/.elysiaclaw:Z
EnvironmentFile={{elysiaclaw_HOME}}/.elysiaclaw/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
